# .config/supervisord.conf
[supervisord]
user = root
logfile = /var/log/supervisor.log

[program:nginx]
command = nginx

[program:gunicorn]
command = gunicorn -c /srv/parents/.config/gunicorn.py mysite.wsgi